# agent/dns_agent.service
[Unit]
Description=DNS Monitor Agent
Documentation=https://github.com/pkthegod/dns-monitor
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=120s

[Service]
Type=simple
User=dns-agent
Group=dns-agent

WorkingDirectory=/opt/dns-agent
ExecStart=/opt/dns-agent/venv/bin/python3 /opt/dns-agent/dns_agent.py

# Segredos injetados via arquivo de ambiente — não aparecem no repositório
# Crie /etc/dns-agent/env com as 3 variáveis abaixo
EnvironmentFile=/etc/dns-agent/env

Restart=on-failure
RestartSec=10s

# Hardening básico
# NoNewPrivileges desabilitado — o agente precisa de sudo para controlar serviços DNS
PrivateTmp=true
ProtectSystem=strict
# /opt/dns-agent precisa ser writable pra o auto-update funcionar:
# o agent escreve dns_agent.py.update_tmp ao lado do binário antes do rename atômico.
# Sem isso, "Update" no painel falha com [Errno 30] Read-only file system.
# /var/lib/dns-agent guarda snapshot de stats DNS pra calcular delta entre amostras.
# Trade-off aceitável: backend autentica via Bearer + checksum SHA-256 validado +
# py_compile antes de gravar. Atacante com root local já tem acesso de qualquer jeito.
ReadWritePaths=/var/log/dns-agent /run/sudo /opt/dns-agent /var/lib/dns-agent

# Logs vão para journald automaticamente
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dns-agent

[Install]
WantedBy=multi-user.target
